# silvus-mock/packaging/systemd/silvus-mock.service
[Unit]
Description=Silvus Radio Mock Emulator
Documentation=https://github.com/silvus-mock
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=silvus-mock
Group=silvus-mock
ExecStart=/usr/local/bin/silvusmock
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
TimeoutStopSec=30
EnvironmentFile=/etc/silvus-mock/env

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/silvus-mock /var/lib/silvus-mock
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Network settings
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=silvus-mock

[Install]
WantedBy=multi-user.target
